# roles/remod-sauerbraten-src/templates/remod.service.j2
[Unit]
Description=Sauerbraten server mod
After=network.target

[Service]
WorkingDirectory={{remod_install_dir}}
ExecStart={{remod_bin}} -q{{remod_install_dir}}

[Install]
WantedBy=multi-user.target
